# deploy/searcher.service
[Unit]
Description=Product Searcher — Gunicorn + Telegram Bot
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/home/product-scanner

EnvironmentFile=-/home/product-scanner/.env
Environment=GUNICORN_BIND=0.0.0.0:8888
Environment=GUNICORN_WORKERS=1
Environment=GUNICORN_THREADS=4
Environment=GUNICORN_TIMEOUT=180

ExecStart=/home/product-scanner/scripts/run-searcher.sh

Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

StandardOutput=append:/home/product-scanner/logs/service.log
StandardError=append:/home/product-scanner/logs/service-error.log

[Install]
WantedBy=multi-user.target
